# systemd/.config/systemd/user/check-battery.timer
[Unit]
Description=Warn if battery is low

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
