# systemd/.config/systemd/user/empty-trash.timer
[Unit]
Description=Trim the FreeDesktop.org Trash bin. 

[Timer]
OnCalendar=Sun,Wed
Persistent=true

[Install]
WantedBy=timers.target
