# systemd/.config/systemd/user/empty-trash.service
[Unit]
Description=Trim the FreeDesktop.org Trash bin. 

[Service]
Type=oneshot
ExecStart=trash-empty 1

[Install]
WantedBy=default.target
