# systemd/.config/systemd/user/i3-session.target
[Unit]
ConditionEnvironment=DISPLAY
Description=i3 User Session
Requires=graphical-session.target
