# systemd/.config/systemd/user/sway-session.target
[Unit]
ConditionEnvironment=WAYLAND_DISPLAY
Description=Sway User Session
Requires=graphical-session.target
