# systemd/.config/systemd/user/sync-email.timer
[Unit]
Description=Copy Gmail contents to local Maildir

[Timer]
OnBootSec=1min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
